# stack/tom/tom.service
[Unit]
Description=T.O.M: Tracking & Observation Microagent
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/tom_agent.py
Restart=always
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
